# templates/nextcloud-backup.timer.j2
[Unit]
Description=NextCloud backup timer

[Timer]
OnCalendar={{ ncbackup_timer }}

[Install]
WantedBy=timers.target
